# kkgroup-api.service
[Unit]
Description=KKCoin Unified API Server
Documentation=https://github.com/chenkankan1103/kkgroup
After=network.target

[Service]
Type=simple
User=e193752468
WorkingDirectory=/home/e193752468/kkgroup
Environment="PYTHONUNBUFFERED=1"
Environment="FLASK_DEBUG=False"
Environment="API_HOST=127.0.0.1"
Environment="API_PORT=5000"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# 啟動命令 - 使用 bash 包裝器
ExecStart=/bin/bash /home/e193752468/kkgroup/start_api.sh

# 重啟策略
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=5

# 日誌設定
StandardOutput=append:/tmp/kkgroup-api.log
StandardError=append:/tmp/kkgroup-api.log

# 安全設定
PrivateTmp=yes
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
