# supervisor_sheet_sync_api.conf
[program:bot-and-api]
# Bot + API 統合進程管理

# 項目路徑（改成 GCP 上的實際路徑，例如 /home/ubuntu/kkgroup）
directory=/path/to/kkgroup

# 執行啟動腳本
command=bash start_bot_and_api.sh

# 自動啟動和重啟
autostart=true
autorestart=true
redirect_stderr=true

# 日誌檔案
stdout_logfile=/var/log/bot-and-api.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

# 運行用戶（改為 GCP 上的用戶名，通常是 ubuntu 或 root）
user=ubuntu

# 環境變量
environment=PATH="/path/to/kkgroup/venv/bin",PYTHONUNBUFFERED=1

# 進程管理
startsecs=10
stopwaitsecs=10
priority=999
