# supervisor_sheet_sync_api.conf
[program:sheet-sync-api]
# Sheet Sync API Gunicorn 配置

# 项目路径
directory=/path/to/kkgroup

# 使用虚拟环境中的 Python 和 gunicorn
command=/path/to/kkgroup/venv/bin/gunicorn -w 4 -b 0.0.0.0:5000 sheet_sync_api:app

# 自动启动和重启
autostart=true
autorestart=true
redirect_stderr=true

# 日志文件
stdout_logfile=/var/log/sheet-sync-api.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

# 运行用户（改为你的用户名，例如 ubuntu）
user=ubuntu

# 环境变量（如需要）
environment=PATH="/path/to/kkgroup/venv/bin",PYTHONUNBUFFERED=1

# 进程管理
startsecs=10
stopwaitsecs=10
