# systemd-backup-20260208-011142/uibot.service
# /etc/systemd/system/discord-uibot.service
[Unit]
Description=Discord UI Bot
After=network.target
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=e193752468
Group=e193752468
WorkingDirectory=/home/e193752468/kkgroup
Environment=PATH=/home/e193752468/kkgroup/venv/bin
ExecStart=/home/e193752468/kkgroup/venv/bin/python uibot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=discord-uibot

TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target
